# vendor/github.com/coreos/coreos-cloudinit/units/user-configdrive.path
[Unit]
Description=Watch for a cloud-config at /media/configdrive

# Note: This unit is essentially just here as a fall-back mechanism to
# trigger cloudinit if it isn't triggered explicitly by other means
# such as by a Wants= in the mount unit. This ensures we handle the
# case where /media/configdrive is provided to a CoreOS container.

[Path]
DirectoryNotEmpty=/media/configdrive
